# launchd/ai.laceyenterprises.adversarial-watcher.placey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.laceyenterprises.adversarial-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/airlock/agent-os/tools/adversarial-review/scripts/adversarial-watcher-start-placey.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/airlock/agent-os/tools/adversarial-review</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/placey</string>
        <key>CODEX_AUTH_PATH</key>
        <string>/Users/placey/.codex/auth.json</string>
        <!-- HQ lives under the airlock-managed Agent OS control-plane root.
             The watcher itself does not consume HQ_ROOT today, but pinning
             it here is the documented contract for higher-trust workers
             (codex / claude-code) that run as placey but spawn remediation
             through HQ-owned infrastructure. See projects/agent-os-trust/
             RUNBOOK-cross-tier-spawn.md for the operating contract. -->
        <key>HQ_ROOT</key>
        <string>/Users/airlock/agent-os-hq</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/placey/Library/Logs/adversarial-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/placey/Library/Logs/adversarial-watcher.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
