# launchd/ai.laceyenterprises.adversarial-follow-up.airlock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Per-user LaunchAgent for the adversarial-review follow-up daemon running
  as airlock. Sibling of ai.laceyenterprises.adversarial-follow-up.placey.plist;
  see that file for the long-lived process model + 1Password / TCC popup
  rationale (the design notes apply equally here).

  Filename convention + auto-install: see the watcher's `.airlock.plist`
  for the full migration rationale. The os-restart path copies this
  file to /Users/airlock/Library/LaunchAgents/ and bootstraps it into
  airlock's gui domain.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.laceyenterprises.adversarial-follow-up.airlock</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/airlock/agent-os/tools/adversarial-review/scripts/adversarial-follow-up-tick.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/airlock/agent-os/tools/adversarial-review</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AGENT_OS_SESSION_LEDGER_DSN</key>
        <string>postgresql://airlock@127.0.0.1:6432/agent_os_ledger</string>
        <key>AGENT_OS_SESSION_LEDGER_TARGET</key>
        <string>postgresql://airlock@127.0.0.1:6432/agent_os_ledger</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/airlock</string>
        <key>USER</key>
        <string>airlock</string>
        <key>LOGNAME</key>
        <string>airlock</string>
        <key>SHELL</key>
        <string>/bin/zsh</string>
        <key>HQ_ROOT</key>
        <string>/Users/airlock/agent-os-hq</string>
        <!-- Enables the final-pass-on-budget-exhausted dispatch lane in
             follow-up-merge-agent.mjs. When the round budget is exhausted
             and the verdict is still `Request changes`, dispatch the
             merge-agent with that trigger so its sub-worker can apply the
             reviewer's findings inline and trigger a fresh review pass.
             See docs/SPEC-adversarial-review-auto-remediation.md for the
             full contract (universal hard gates + blocker-class refusal
             still apply). -->
        <key>MERGE_AGENT_FINAL_PASS_ON_REQUEST_CHANGES</key>
        <string>1</string>

        <!--
          ALERT_TO 1Password ref is REQUIRED. See the corresponding watcher
          plist for the resolution contract. Mirrored here because the
          follow-up daemon's tick path validates the same env at startup.
        -->
        <key>ADVERSARIAL_REVIEW_ALERT_TO_OP_REF</key>
        <string>op://Cliovault/adversarial-watcher-alert-to/credential</string>

        <!--
          BUILD-PACK-DEFAULT-FLIP / 2026-05-29 → 2026-05-30 → reverted
          2026-05-30 (PT) when codex OAuth came back online. Per the
          operator directive: remediator pins to `codex` as an explicit
          launchd override. The module config.yaml keeps its default for
          non-launchd/dev runs.

          2026-06-18: reviewer env pin REMOVED here to match the watcher
          LaunchAgent (#333). The reviewer role now resolves from config —
          roles.reviewer: adversarial in config.local.yaml — restoring the
          cross-model claude<->codex builder-never-reviews routing. This
          env pin (Layer 5) was overriding the yaml (Layer 4) and is
          validated by the follow-up daemon's tick path, so leaving it here
          would force a Codex reviewer for follow-up runs and waive the
          cross-model guarantee for same-family builders even after the
          watcher flipped back. The remediator pin stays `codex` (unchanged).
          This is an explicit launchd override, not cross-model-by-default.

          See watcher plist for full rationale + LAC-893 and the role-pin
          env knobs documented in `docs/follow-up-runbook.md`.
        -->
        <key>ADVERSARIAL_REVIEW_DEFAULT_REMEDIATOR</key>
        <string>claude-code</string>
        <!--
          Match the watcher LaunchAgent: follow-up comment posting uses the
          same reviewer identities and should mint reviewer GitHub tokens
          through the local OAuth broker instead of per-role 1Password PAT
          reads.
        -->
        <key>CLAUDE_REVIEWER_AUTH_VIA_BROKER</key>
        <string>true</string>
        <key>CODEX_REVIEWER_AUTH_VIA_BROKER</key>
        <string>true</string>
        <key>GEMINI_REVIEWER_AUTH_VIA_BROKER</key>
        <string>true</string>
        <key>OAUTH_BROKER_SHARED_SECRET_FILE</key>
        <string>/Users/airlock/agent-os/.secrets/local/oauth-broker-shared-secret</string>
        <key>OAUTH_BROKER_CLAUDE_REVIEWER_EXPECTED_APP_ID</key>
        <string>3994110</string>
        <key>OAUTH_BROKER_CLAUDE_REVIEWER_EXPECTED_INSTALLATION_ID</key>
        <string>138772505</string>
        <key>OAUTH_BROKER_CODEX_REVIEWER_EXPECTED_APP_ID</key>
        <string>3994060</string>
        <key>OAUTH_BROKER_CODEX_REVIEWER_EXPECTED_INSTALLATION_ID</key>
        <string>138771486</string>
        <key>OAUTH_BROKER_GEMINI_REVIEWER_EXPECTED_APP_ID</key>
        <string>3994139</string>
        <key>OAUTH_BROKER_GEMINI_REVIEWER_EXPECTED_INSTALLATION_ID</key>
        <string>138773049</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/airlock/Library/Logs/adversarial-follow-up.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/airlock/Library/Logs/adversarial-follow-up.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ExitTimeOut</key>
    <integer>45</integer>
</dict>
</plist>
